# smart-pause-resume.service
[Unit]
Description=Auto-pause previous media when new media starts

[Service]
ExecStart=%h/bin/smart-pause-resume
Restart=on-failure

[Install]
WantedBy=default.target
